# systemd/themenizer.service.template
[Unit]
Description=Run themenizer script daily at 6 AM
After=graphical-session.target

[Install]
WantedBy=graphical-session.target

[Service]
